# deploy/systemd/whrkhldsb-ssh-ws.service.example
[Unit]
Description={{SITE_NAME}} SSH WebSocket proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{APP_DIR}}
EnvironmentFile={{ENV_FILE}}
Environment=NODE_ENV=production
Environment=PATH={{SYSTEMD_PATH}}
ExecStart={{NPX_BIN}} tsx src/ssh-ws-proxy.ts
Restart=always
RestartSec=5
User={{APP_USER}}
Group={{APP_USER}}

[Install]
WantedBy=multi-user.target
